# ansible/collections/ansible_collections/labo/infra/roles/local_ansible_server/files/lansible-playbook.service
[Unit]
Description=/opt/ansible-nonroot/bin/_lansible-playbook

[Service]
Type=simple
ExecStart=/opt/ansible-nonroot/bin/_lansible-playbook
